# mural.service
[Unit]
Description=Mural

[Service]
User=pi
ExecStart=/usr/bin/python3 /home/pi/MuralBot/main.py

[Install]
WantedBy=multi-user.target
